# src/units/user/dbus-broker.service.in
[Unit]
Description=D-Bus User Message Bus
Documentation=man:dbus-broker-launch(1)
DefaultDependencies=false

After=dbus.socket dbus-metrics.socket
Before=basic.target shutdown.target
Conflicts=shutdown.target
Requires=dbus.socket
Wants=dbus-metrics.socket

[Service]
Type=notify-reload
ExecStart=@bindir@/dbus-broker-launch --scope user
Slice=session.slice

[Install]
Alias=dbus.service
