# units/user/dbus-broker.service.in
[Unit]
Description=D-Bus User Message Bus
Documentation=man:dbus-broker(1)
Requires=dbus.socket

[Service]
ExecStart=@bindir@/dbus-broker-launch -v --scope user --listen inherit
