# units/user/dbus-broker.service.in
[Unit]
Description=D-Bus User Message Bus
Documentation=man:dbus-broker-launch(1)
DefaultDependencies=false
Before=basic.target shutdown.target
Conflicts=shutdown.target

[Service]
Sockets=dbus.socket
ExecStart=@bindir@/dbus-broker-launch -v --scope user --listen inherit
ExecReload=@bindir@/busctl call org.freedesktop.DBus /org/freedesktop/DBus org.freedesktop.DBus ReloadConfig

[Install]
Alias=dbus.service
