# units/system/dbus-broker.service.in
[Unit]
Description=D-Bus System Message Bus
Documentation=man:dbus-broker(1)
Requires=dbus.socket

[Service]
OOMScoreAdjust=-900
LimitNOFILE=16384
ExecStart=@bindir@/dbus-broker-launch -v --scope system --listen inherit

[Install]
Alias=dbus.service
